# systemd/system-report.service
[Unit]
Description=Run the system report script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/system_report.sh
